# deploy/systemd/sharon-kitchen.service.example
[Unit]
Description=Sharon Kitchen Travel website
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/sharon-kitchen-travel
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=TRUST_PROXY=1

[Install]
WantedBy=multi-user.target
